# rc.d/vpnsw.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: vpnsw
# REQUIRE: LOGIN

. /etc/rc.subr

name="vpnsw"
rcvar="vpnsw_enable"

pidfile="@APP_RUNDIR@"/vpnsw.pid
command="@prefix@/sbin/${name}"
command_args=""
procname="@prefix@/sbin/${name}"
command_interpreter="@PERL@"

load_rc_config ${name}

: ${pgagent_enable:="NO"}

run_rc_command "$1"
#EOF
