# vpnsw.service.in
[Unit]
Description=vpnsw

[Service]
Type=forking
PIDFile=@APP_RUNDIR@/vpnsw.pid
ExecStart=@prefix@/sbin/vpnsw
ExecReload=/bin/kill -HUP $MAINPID
RuntimeDirectory=@APP_RUNDIR@
ExecStartPre=/bin/mkdir -p @APP_RUNDIR@

[Install]
WantedBy=multi-user.target

